# skills/box-rclone-binder/scripts/install/box-binder-health.timer.tmpl
[Unit]
Description=box-binder health timer

[Timer]
OnCalendar=*:0/15
Persistent=true
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
